# pkg-plist
bin/%%USER%%/agent
bin/%%USER%%/dist/checks/__init__.py
bin/%%USER%%/dist/checks/libs/__init__.py
bin/%%USER%%/dist/checks/libs/wmi/__init__.py
bin/%%USER%%/dist/checks/libs/wmi/sampler.py
bin/%%USER%%/dist/checks/network_checks.py
bin/%%USER%%/dist/checks/prometheus_check/__init__.py
bin/%%USER%%/dist/checks/winwmi_check.py
bin/%%USER%%/dist/conf.d/apm.yaml.default
bin/%%USER%%/dist/conf.d/containerd.d/conf.yaml.example
bin/%%USER%%/dist/conf.d/cpu.d/conf.yaml.default
bin/%%USER%%/dist/conf.d/cri.d/conf.yaml.example
bin/%%USER%%/dist/conf.d/disk.d/conf.yaml.default
bin/%%USER%%/dist/conf.d/docker.d/conf.yaml.example
bin/%%USER%%/dist/conf.d/file_handle.d/conf.yaml.default
bin/%%USER%%/dist/conf.d/go_expvar.d/agent_stats.yaml.example
bin/%%USER%%/dist/conf.d/io.d/conf.yaml.default
bin/%%USER%%/dist/conf.d/jmx.d/conf.yaml.example
bin/%%USER%%/dist/conf.d/kubernetes_apiserver.d/conf.yaml.example
bin/%%USER%%/dist/conf.d/load.d/conf.yaml.default
bin/%%USER%%/dist/conf.d/memory.d/conf.yaml.default
bin/%%USER%%/dist/conf.d/network.d/conf.yaml.default
bin/%%USER%%/dist/conf.d/ntp.d/conf.yaml.default
bin/%%USER%%/dist/conf.d/oom_kill.d/conf.yaml.example
bin/%%USER%%/dist/conf.d/process_agent.yaml.default
bin/%%USER%%/dist/conf.d/systemd.d/conf.yaml.example
bin/%%USER%%/dist/conf.d/tcp_queue_length.d/conf.yaml.example
bin/%%USER%%/dist/conf.d/uptime.d/conf.yaml.default
bin/%%USER%%/dist/conf.d/winproc.d/conf.yaml.default
bin/%%USER%%/dist/config.py
bin/%%USER%%/dist/%%USER%%.yaml
bin/%%USER%%/dist/dd-agent
bin/%%USER%%/dist/system-probe.yaml
bin/%%USER%%/dist/templates/aggregator.tmpl
bin/%%USER%%/dist/templates/clusteragent.tmpl
bin/%%USER%%/dist/templates/collector.tmpl
bin/%%USER%%/dist/templates/compliance.tmpl
bin/%%USER%%/dist/templates/custommetricsprovider.tmpl
bin/%%USER%%/dist/templates/dogstatsd.tmpl
bin/%%USER%%/dist/templates/endpoints.tmpl
bin/%%USER%%/dist/templates/forwarder.tmpl
bin/%%USER%%/dist/templates/header.tmpl
bin/%%USER%%/dist/templates/jmxfetch.tmpl
bin/%%USER%%/dist/templates/logsagent.tmpl
bin/%%USER%%/dist/templates/metadatamapper.tmpl
bin/%%USER%%/dist/templates/runtimesecurity.tmpl
bin/%%USER%%/dist/templates/snmp-traps.tmpl
bin/%%USER%%/dist/templates/systemprobe.tmpl
bin/%%USER%%/dist/templates/trace-agent.tmpl
bin/%%USER%%/dist/utils/__init__.py
bin/%%USER%%/dist/utils/containers.py
bin/%%USER%%/dist/utils/platform.py
bin/%%USER%%/dist/utils/py_mem.py
bin/%%USER%%/dist/utils/subprocess_output.py
bin/%%USER%%/dist/utils/tailfile.py
bin/%%USER%%/dist/utils/timeout.py
%%PROCESS%%bin/%%USER%%/process-agent
%%APM%%bin/%%USER%%/trace-agent
%%ETCDIR%%/%%USER%%.yaml.example
etc/rc.d/%%USER%%-agent
etc/rc.d/%%USER%%-process-agent
etc/rc.d/%%USER%%-trace-agent
include/%%USER%%_agent_rtloader.h
include/rtloader_mem.h
include/rtloader_types.h
lib/lib%%USER%%-agent-rtloader.so
lib/lib%%USER%%-agent-rtloader.so.0.1.0
lib/lib%%USER%%-agent-rtloader.so.1
lib/lib%%USER%%-agent-three.so
%%PORTDOCS%%%%DOCSDIR%%/CHANGELOG.rst
%%PORTDOCS%%%%DOCSDIR%%/CONTRIBUTING.md
%%PORTDOCS%%%%DOCSDIR%%/LICENSE
%%PORTDOCS%%%%DOCSDIR%%/README.md
@dir %%ETCDIR%%/conf.d
@dir(%%USER%%,%%GROUP%%,0750) %%LOGDIR%%
@dir(%%USER%%,%%GROUP%%,0750) %%RUNDIR%%
