# files/supervisord.conf.in
[supervisord]
logfile = %%LOGDIR%%/supervisord.log
logfile_maxbytes = 50MB
loglevel = warn
nodaemon = false
identifier = supervisor
nocleanup = true
pidfile = %%PIDDIR%%/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = %%PIDDIR%%/agent-supervisor.sock

[supervisorctl]
prompt = datadog
serverurl = unix://%%PIDDIR%%/agent-supervisor.sock

[program:collector]
command=%%PYTHON_CMD%% %%PYTHON_SITELIBDIR%%/datadog/agent.py foreground --use-local-forwarder
process_name=%(program_name)s
stdout_logfile=%%LOGDIR%%/collector.log
redirect_stderr=true
priority=999
startsecs=2
environment=LANG=POSIX,PYTHONPATH='datadog/checks/libs:$PYTHONPATH'

[program:forwarder]
command=%%PYTHON_CMD%% %%PYTHON_SITELIBDIR%%/datadog/ddagent.py --use_simple_http_client=1
process_name=%(program_name)s
stdout_logfile=%%LOGDIR%%/forwarder.log
redirect_stderr=true
priority=998
startsecs=3

[program:dogstatsd]
command=%%PYTHON_CMD%% %%PYTHON_SITELIBDIR%%/datadog/dogstatsd.py --use-local-forwarder
process_name=%(program_name)s
stdout_logfile=%%LOGDIR%%/dogstatsd.log
redirect_stderr=true
priority=998
startsecs=3

[program:jmxfetch]
command=%%PYTHON_CMD%% %%PYTHON_SITELIBDIR%%/datadog/jmxfetch.py
stdout_logfile=%%LOGDIR%%/jmxfetch.log
redirect_stderr=true
priority=999
startsecs=3

[group:datadog-agent]
programs=forwarder,collector,dogstatsd,jmxfetch

