# files/supervisord.conf.in
[supervisord]
logfile = %%LOGDIR%%/supervisord.log
logfile_maxbytes = 50MB
loglevel = info
nodaemon = true
identifier = supervisor
nocleanup = true
pidfile = %%PIDDIR%%/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = %%PIDDIR%%/agent-supervisor.sock

[supervisorctl]
prompt = datadog
serverurl = unix://%%PIDDIR%%/agent-supervisor.sock

[program:collector]
command=python %%PYTHON_SITELIBDIR%%/dd-agent/agent.py foreground --use-local-forwarder
stdout_logfile=%%LOGDIR%%/collector.log
redirect_stderr=true
priority=999
startsecs=2
environment=LANG=POSIX,PYTHONPATH='agent/checks/libs:$PYTHONPATH'

[program:forwarder]
command=python %%PYTHON_SITELIBDIR%%/dd-agent/ddagent.py --use_simple_http_client=1
stdout_logfile=%%LOGDIR%%/forwarder.log
redirect_stderr=true
priority=998
startsecs=3

[program:dogstatsd]
command=python %%PYTHON_SITELIBDIR%%/dd-agent/dogstatsd.py --use-local-forwarder
stdout_logfile=%%LOGDIR%%/dogstatsd.log
redirect_stderr=true
priority=998
startsecs=3

[program:jmxfetch]
command=python %%PYTHON_SITELIBDIR%%/dd-agent/jmxfetch.py
stdout_logfile=%%LOGDIR%%/jmxfetch.log
redirect_stderr=true
priority=999
startsecs=3

[group:datadog-agent]
programs=forwarder,collector,dogstatsd,jmxfetch

